# templates/systemd/certbot-renew.timer
[Unit]
Description=Timer for certbot renewal

[Timer]
OnBootSec=300
OnUnitActiveSec=1w

[Install]
WantedBy=multi-user.target